# systemd/user/goimapnotify@.service.d/override.conf
[Unit]
PartOf=mail.target
BindsTo=davmail.service
After=mail.target davmail.service

[Install]
WantedBy=mail.target
