# systemd/user/goimapnotify@.service.d/override.conf
[Unit]
PartOf=mail.target
After=mail.target

[Install]
WantedBy=mail.target
