# systemd/user/goimapnotify@.service.d/override.conf
[Unit]
PartOf=mail.target
After=mail.target

[Service]
Environment="DISPLAY=:0"

[Install]
WantedBy=mail.target
