# systemd/user/goimapnotify@.service.d/override.conf
[Unit]
Requires=mail.target davmail.service
After=davmail.service

[Install]
WantedBy=mail.target