# systemd/user/vdirsyncer.timer.d/override.conf
[Unit]
Requires=mail.target davmail.service
PartOf=mail.target
BindsTo=davmail.service
After=mail.target davmail.service
