# systemd/user/vdirsyncer.timer.d/override.conf
[Unit]
Requires=mail.target davmail.service
