# systemd/user/dunst.service
[Unit]
Description=Dunst User Service

[Service]
Environment="DISPLAY=:0"
ExecStart=/usr/bin/dunst

[Install]
WantedBy=multi-user.target
