# systemd/user/mail.target
[Unit]
Description=Mail target to manage dependent mail units
AllowIsolate=no
Wants=davmail.service

[Install]
WantedBy=default.target
