# systemd/user/unclutter.service
[Unit]
Description="Unclutters the screen by removing the mouse"

[Service]
ExecStart=/usr/bin/unclutter

[Install]
WantedBy=default.target
